# ids.service
[Unit]
    Description= Mon petit service

    [Service]
    ExecStart= *programme a lancer*